# pybeautifier.service
[Unit]
Description=PyBeautifier daemon

[Service]
Type=forking
ExecStart=/usr/bin/pybeautifier -d

[Install]
WantedBy=multi-user.target
